# examples/05_graceful_restart/setup/supervisord.conf
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/run/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:helloWorker]
autorestart = true ;the process will be unconditionally restarted when it exits, without regard to its exit code
command = php /home/dbaltas/projects/rnd/gearman-php-prototype/examples/05_graceful_restart/helloWorker.php
numprocs = 3                    ; number of processes copies to start (def 1)
process_name = %(program_name)s%(process_num)s ; process_name expr (default %(program_name)s)
stdout_logfile= /var/log/workers/stdout.log
stderr_logfile= /tmp/error.log
startsecs=0
autostart=true
autorestart=true
user=root